# packages/systemd/user/prettierd.service
[Unit]
Description=Prettierd daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/prettierd
Restart=always
RestartSec=3
Environment=PRETTIERD_DEFAULT_CONFIG=%h/.config/prettierd/.prettierrc

[Install]
WantedBy=default.target
