# packages/systemd/user/tmux-wiki-session.service
[Unit]
Description=Start tmux wiki session

[Service]
Type=oneshot
ExecStart=/usr/bin/tmux new-session -d -s wiki -c /mnt/share/wiki
RemainAfterExit=yes

[Install]
WantedBy=default.target
