# packages/litellm/system/runit/etc/sv/litellm/run
#!/bin/sh
exec uv tool run litellm proxy --config /home/shad/.config/litellm/config.yaml
