# packages/litellm/system/systemd/etc/systemd/user/litellm.service
[Unit]
Description=LiteLLM Proxy Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/env uv tool run litellm --config %h/.config/litellm/config.yaml
Restart=on-failure
RestartSec=5s
Environment=MINIMAX_API_KEY
Environment=LITELLM_MASTER_KEY

[Install]
WantedBy=default.target